# PlatformAgents/com/cognizant/devops/platformagents/agents/ci/citfs/InSightsCITFSAgent.service
[Unit]
Description=Ubuntu CITFSAgent Service

[Service]
Type=forking
ExecStart=/bin/bash /opt/insightsagents/PlatformAgents/citfs/InSightsCITFSAgent.sh

[Install]
WantedBy=multi-user.target