# PlatformAgents/com/cognizant/devops/platformagents/agents/deployment/xldeploy/InSightsXLDeployAgent.service
[Unit]
Description=Ubuntu XLDeployAgent Service

[Service]
Type=forking
ExecStart=/bin/bash /opt/insightsagents/PlatformAgents/xldeploy/InSightsXLDeployAgent.sh

[Install]
WantedBy=multi-user.target