# PlatformAgents/com/cognizant/devops/platformagents/agents/ci/concourse/InSightsConcourseAgent.service
[Unit]
Description=Ubuntu ConcourseAgent Service

[Service]
Type=forking
ExecStart=/bin/bash /opt/insightsagents/PlatformAgents/concourse/InSightsConcourseAgent.sh

[Install]
WantedBy=multi-user.target