# PlatformAgents/com/cognizant/devops/platformagents/agents/deployment/buildmaster/buildmasteragent.service
[Unit]
Description=Ubuntu buildmaster Service

[Service]
Type=forking
ExecStart=/bin/bash /opt/insightsagents/PlatformAgents/buildmaster/__AGENT_KEY__/__AGENT_KEY__.sh

[Install]
WantedBy=multi-user.target